# package/hotspot-trial/files/etc/uci-defaults/54-hotspot-trial
#!/bin/sh
# Add cron job for daily trial reset
if ! grep -q "trial-reset.sh" /etc/crontabs/root 2>/dev/null; then
    echo "0 0 * * * /usr/lib/hotspotos/trial-reset.sh reset" >> /etc/crontabs/root
    /etc/init.d/cron restart 2>/dev/null || true
fi
exit 0
