# package/hotspot-server/files/etc/uci-defaults/52-hotspot-server
#!/bin/sh
/etc/init.d/hotspot-server enable 2>/dev/null
exit 0
