# package/hotspot-users/files/etc/uci-defaults/53-hotspot-users
#!/bin/sh
# Initialize default admin user in database
/usr/lib/hotspotos/users.sh init-defaults 2>/dev/null
exit 0
